# utils/temperature-to-cloud.service.in
[Unit]
Description=Temperature monitor - send to cloud
Wants=temperatore-server.service
After=network.target

[Service]
ExecStart=/usr/bin/python3 -u temperature-to-cloud.py
WorkingDirectory=@CMAKE_INSTALL_PREFIX@
StandardOutput=append:/var/log/temperature-monitor.log
StandardError=append:/var/log/temperature-monitor.log
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
