# utils/temperature-server.service.in
[Unit]
Description=Temperature monitor server
After=network.target

[Service]
ExecStart=temperature-server
WorkingDirectory=@CMAKE_INSTALL_PREFIX@/bin
StandardOutput=append:/var/log/temperature-monitor.log
StandardError=append:/var/log/temperature-monitor.log
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
