# utils/temperature-server.service.in
[Unit]
Description=Temperature monitor server
After=network.target

[Service]
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/temperature-server
#StandardOutput=append:/var/log/temperature-monitor.log
#StandardError=append:/var/log/temperature-monitor.log
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
